# infra/saista-payment.service
[Unit]
Description=Saista Bakers - Payment Service (FastAPI :5003)
After=network.target saista-user.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/saista/payment-service
EnvironmentFile=/opt/saista/payment-service/.env
ExecStart=/opt/saista/payment-service/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 5003 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target